# private_Library/LaunchAgents/gitwatch.org.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>gitwatch.org</string>
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/gitwatch/gitwatch-org-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/gitwatch/gitwatch-org-stderr.log</string>
	<key>LimitLoadToSessionType</key>
	<array>
	  <string>Aqua</string>
	  <string>Background</string>
	  <string>LoginWindow</string>
	  <string>StandardIO</string>
	  <string>System</string>
	</array>
	<key>ProgramArguments</key>
	<array>
      <string>{{ .path.brew }}/gitwatch</string>
      <string>-r</string>
      <string>origin</string>
      <string>{{ .chezmoi.homeDir }}/org</string>
	</array>
	<key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{ .path.brew }}:/bin:/usr/bin</string>
    </dict>
	<key>RunAtLoad</key>
	<true/>
  </dict>
</plist>
